# contrib/init/sallmond.openrc
#!/sbin/runscript

# backward compatibility for existing gentoo layout 
#
if [ -d "/var/lib/sallmon/.sallmon" ]; then
	SALLMOND_DEFAULT_DATADIR="/var/lib/sallmon/.sallmon"
else
	SALLMOND_DEFAULT_DATADIR="/var/lib/sallmond"
fi

SALLMOND_CONFIGFILE=${SALLMOND_CONFIGFILE:-/etc/sallmon/sallmon.conf}
SALLMOND_PIDDIR=${SALLMOND_PIDDIR:-/var/run/sallmond}
SALLMOND_PIDFILE=${SALLMOND_PIDFILE:-${SALLMOND_PIDDIR}/sallmond.pid}
SALLMOND_DATADIR=${SALLMOND_DATADIR:-${SALLMOND_DEFAULT_DATADIR}}
SALLMOND_USER=${SALLMOND_USER:-${SALLMON_USER:-sallmon}}
SALLMOND_GROUP=${SALLMOND_GROUP:-sallmon}
SALLMOND_BIN=${SALLMOND_BIN:-/usr/bin/sallmond}
SALLMOND_NICE=${SALLMOND_NICE:-${NICELEVEL:-0}}
SALLMOND_OPTS="${SALLMOND_OPTS:-${SALLMON_OPTS}}"

name="Sallmon Core Daemon"
description="Sallmon cryptocurrency P2P network daemon"

command="/usr/bin/sallmond"
command_args="-pid=\"${SALLMOND_PIDFILE}\" \
		-conf=\"${SALLMOND_CONFIGFILE}\" \
		-datadir=\"${SALLMOND_DATADIR}\" \
		-daemon \
		${SALLMOND_OPTS}"

required_files="${SALLMOND_CONFIGFILE}"
start_stop_daemon_args="-u ${SALLMOND_USER} \
			-N ${SALLMOND_NICE} -w 2000"
pidfile="${SALLMOND_PIDFILE}"

# The retry schedule to use when stopping the daemon. Could be either
# a timeout in seconds or multiple signal/timeout pairs (like
# "SIGKILL/180 SIGTERM/300")
retry="${SALLMOND_SIGTERM_TIMEOUT}"

depend() {
	need localmount net
}

# verify
# 1) that the datadir exists and is writable (or create it)
# 2) that a directory for the pid exists and is writable
# 3) ownership and permissions on the config file
start_pre() {
	checkpath \
	-d \
	--mode 0750 \
	--owner "${SALLMOND_USER}:${SALLMOND_GROUP}" \
	"${SALLMOND_DATADIR}"

	checkpath \
	-d \
	--mode 0755 \
	--owner "${SALLMOND_USER}:${SALLMOND_GROUP}" \
	"${SALLMOND_PIDDIR}"

	checkpath -f \
	-o ${SALLMOND_USER}:${SALLMOND_GROUP} \
	-m 0660 \
	${SALLMOND_CONFIGFILE}

	checkconfig || return 1
}

checkconfig()
{
	if ! grep -qs '^rpcpassword=' "${SALLMOND_CONFIGFILE}" ; then
		eerror ""
		eerror "ERROR: You must set a secure rpcpassword to run sallmond."
		eerror "The setting must appear in ${SALLMOND_CONFIGFILE}"
		eerror ""
		eerror "This password is security critical to securing wallets "
		eerror "and must not be the same as the rpcuser setting."
		eerror "You can generate a suitable random password using the following "
		eerror "command from the shell:"
		eerror ""
		eerror "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
		eerror ""
		eerror "It is recommended that you also set alertnotify so you are "
		eerror "notified of problems:"
		eerror ""
		eerror "ie: alertnotify=echo %%s | mail -s \"Sallmon Alert\"" \
			"admin@foo.com"
		eerror ""
		return 1
	fi
}
